# ansible/nishantnikhil-github/roles/jenkins/templates/jenkins.service
[Unit]
Description=Daemon Jenkins2

[Service]
ExecStart=/usr/bin/java -jar /usr/local/etc/jenkins.war

[Install]
WantedBy=multi-user.target
